# roles/infopi/templates/chromium-tab-switch.timer.j2
# {{ ansible_managed }}

[Unit]
Description=Run chromium-tab-switcher periodically

[Timer]
OnCalendar={{ infopi_tab_timer }}
AccuracySec=1sec
Persistent=true

[Install]
WantedBy=timers.target
